# rootfs/etc/services.d/bus/run
#!/usr/bin/with-contenv sh
java -jar /app/karate-0.9.1.jar -m /app/bus-mock.feature -p 8087
